# supervisor/cozy-controller.conf
[program:cozy-controller]
autorestart=true
command=node server.js
directory=/usr/local/lib/node_modules/cozy-controller/build
redirect_stderr=true
user=root